# etc/supervisor/etcd.conf
#this for test configuration with etcd
[program:etcd]
command=etcd
user=root
startsecs=1
autostart=true
autorestart=true
stdout_events_enabled = true
stderr_events_enabled = true